# server-config/supervisor/conf.d/cron.conf
[program:cron]
command=/usr/sbin/cron
user=root
autostart=true